# targets/slideruleearth/docker/monitor-agent/rootfs/etc/services.d/node_exporter/finish
#!/bin/sh
echo "node_exporter exited with code $1"