# targets/slideruleearth/docker/monitor/rootfs/etc/services.d/loki/run
#!/command/execlineb -P
loki -config.file=/loki/loki.yml