# targets/slideruleearth/docker/monitor-agent/rootfs/etc/services.d/promtail/finish
#!/bin/sh
echo "promtail exited with code $1"