# targets/slideruleearth/docker/monitor/rootfs/etc/services.d/prometheus/finish
#!/bin/sh
echo "grafana exited with code $1"