# targets/slideruleearth/docker/monitor/rootfs/etc/services.d/grafana/run
#!/command/execlineb -P
grafana-server \
    --homepath "/grafana" \
    --config "/grafana/conf/grafana.ini" \
    web